# root/etc/s6-overlay/s6-rc.d/svc-mod-qbittorrent-freezetag/run
#!/usr/bin/with-contenv bash
LOGFILE="/tmp/freezetag.log"

echo "Starting svc-qbittorrent-freezetag..." >> $LOGFILE
date >> $LOGFILE

# This is an example service that would run for the mod
# It depends on init-services, the baseimage hook for start of all longrun services


exec \
    freezetag mount \
    --uid "${PUID}" \
    --gid "${PGID}" \
    "/data/torrents" \
    "/data/torrents_frozen"

echo "Dismounting mod-qbittorrent-freezetag service" >> $LOGFILE 2>&1