# root/etc/s6-overlay/s6-rc.d/svc-mod-qbittorrent-freezetag/run
#!/usr/bin/with-contenv bash

# This is an example service that would run for the mod
# It depends on init-services, the baseimage hook for start of all longrun services

echo "Starting mod-qbittorrent-freezetag service (mount)"

exec \
    freezetag mount
        --uid "${PUID}"
        --gid "${PGID}"
        "/data/torrents" "/data/torrents_frozen"

echo "Dismounting mod-qbittorrent-freezetag service"