# root/etc/s6-overlay/s6-rc.d/svc-mod-qbittorrent-freezetag/run
#!/usr/bin/with-contenv bash

# Check if DOCKER_MODS_DEBUG is set and true
if [[ "${DOCKER_MODS_DEBUG}" == "true" ]]; then

    log() {
        echo "[mod-freezetag] - $1"
    }

    # Initial logging
    log "Starting svc-qbittorrent-freezetag..."
else
    log() {
        # No-op function to avoid logging when DOCKER_MODS_DEBUG is not true
        :
    }
fi



# Check if the CATEGORY_DIR and FROZEN_DIR environment variables are set and that the directories exist
if [[ -z "${CATEGORY_DIR}" ]]; then
    log "Error: The CATEGORY_DIR environment variables must be set."
    exit 1
fi

if [[ -z "${FROZEN_DIR}" ]]; then
    log "Error: The FROZEN_DIR environment variables must be set."
    exit 1
fi

if [[ ! -d "${CATEGORY_DIR}" ]]; then
    log "Error: The CATEGORY_DIR (${CATEGORY_DIR}) does not exist."
    exit 1
fi

if [[ ! -d "${FROZEN_DIR}" ]]; then
    log "Error: The FROZEN_DIR (${FROZEN_DIR}) does not exist."
    exit 1
fi

# Check if the FROZEN_DIR is empty
if [ ! -z "$(find ${FROZEN_DIR} -mindepth 1 -print -quit)" ]; then
    log "Error: The FROZEN_DIR (${FROZEN_DIR}) is not empty."
    exit 1
fi

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "mountpoint -q ${FROZEN_DIR}" \
        freezetag mount \
        --uid "${PUID}" \
        --gid "${PGID}" \
        "${CATEGORY_DIR}" \
        "${FROZEN_DIR}"

# Log the execution timestamp
log "Service script started."