# root/etc/s6-overlay/s6-rc.d/svc-mod-qbittorrent-freezetag/run
#!/usr/bin/with-contenv bash
# echo "Dummy service for testing started" >> /tmp/dummy_service.log
# while true; do
#     echo "Dummy service is running at $(date)" >> /tmp/dummy_service.log
#     sleep 10
# done

# This is an example service that would run for the mod
# It depends on init-services, the baseimage hook for start of all longrun services

echo "Starting mod-qbittorrent-freezetag service (mount)" > /dev/console 2>&1

exec \
    freezetag mount \
    --uid "${PUID}" \
    --gid "${PGID}" \
    "/data/torrents" \
    "/data/torrents_frozen"

echo "Dismounting mod-qbittorrent-freezetag service" > /dev/console 2>&1