# root/etc/s6-overlay/s6-rc.d/init-mod-qbittorrent-freezetag-postinstall/run
pip3 install git+https://github.com/undaunt/freezetag --break-system-packages